# anywhere_door_agent/deploy/linux/anywhere-door-agent.service
[Unit]
Description=Anywhere Door Agent - OS Level File Watcher (Rust)
After=network.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/anywhere_door_agent
WorkingDirectory=/var/lib/anywhere-door-agent

# File watcher only - memory efficient event streaming
Environment=ANYWHERE_DOOR_ENABLE_OS_WATCHER=true
Environment=ANYWHERE_DOOR_FILE_EVENT_METADATA_OUTPUT=/var/log/anywhere-door-agent/file_event_metadata.ndjson
Environment=ANYWHERE_DOOR_WATCH_ROOTS=/home/deilsy,/var/log/anywhere-door-agent

Restart=always
RestartSec=3
User=anywhere-door
Group=anywhere-door
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target